# scripts/ai.reslu.aria-heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Step 12 aria_queue heartbeat schedule (Mac mini).

  Checks aria_queue for pending items every 5 minutes. If any are
  found, injects a system event via `openclaw system event --mode now`
  to wake Aria immediately. Zero pending = exit 0, zero cost.

  Install:
    cp scripts/ai.reslu.aria-heartbeat.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/ai.reslu.aria-heartbeat.plist
    launchctl kickstart -p gui/$(id -u)/ai.reslu.aria-heartbeat   # optional: run once now
  Disable:
    launchctl bootout gui/$(id -u)/ai.reslu.aria-heartbeat
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.reslu.aria-heartbeat</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/vale/reslu-spec-system/scripts/aria_heartbeat.py</string>
    </array>

    <!-- 5-minute cadence — cheap count check costs zero tokens when queue is empty. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/Users/vale/.openclaw/workspace/vault/agent-openclaw/daily/aria-heartbeat.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vale/.openclaw/workspace/vault/agent-openclaw/daily/aria-heartbeat-errors.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- openclaw lives in /opt/homebrew/bin -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/vale</string>
    </dict>
</dict>
</plist>
