# scripts/ai.reslu.aria-heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Step 12 aria_queue heartbeat schedule (Mac mini).

  Frequent by design — the whole point of a heartbeat is a cheap,
  regular check that costs nothing when idle. StartInterval=120 (2
  min) is a reasonable starting cadence; tighten or loosen once
  wake_aria() is actually finished and you can see real behaviour.
  Same ai.reslu.* LaunchAgent convention as
  ai.reslu.email-ingest.plist (Step 8).

  Uses plain python3 (stdlib only, no venv needed — see
  aria_heartbeat.py's own header) rather than the
  .venv-email-ingest interpreter that script needs for its
  talon/pdftotext dependencies.

  Install (does NOT run until you load it — validate with a manual
  `python3 scripts/aria_heartbeat.py` run first, and finish
  wake_aria() before enabling on a schedule):
    cp scripts/ai.reslu.aria-heartbeat.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/ai.reslu.aria-heartbeat.plist
    launchctl kickstart -p gui/$(id -u)/ai.reslu.aria-heartbeat   # optional: run once now
  Disable:
    launchctl bootout gui/$(id -u)/ai.reslu.aria-heartbeat
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.reslu.aria-heartbeat</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/vale/reslu-spec-system/scripts/aria_heartbeat.py</string>
    </array>

    <!-- 2-minute cadence — cheap enough to run often, see header note. -->
    <key>StartInterval</key>
    <integer>120</integer>

    <key>StandardOutPath</key>
    <string>/Users/vale/.openclaw/workspace/vault/agent-openclaw/daily/aria-heartbeat.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vale/.openclaw/workspace/vault/agent-openclaw/daily/aria-heartbeat-errors.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/vale</string>
    </dict>
</dict>
</plist>
